# COPY/etc/systemd/system/rh-postgresql95-repmgr.service.d/config_file.conf
[Service]
PIDFile=/var/opt/rh/rh-postgresql95/run/repmgr/repmgrd-9.5.pid

Environment=REPMGRDCONF=/etc/repmgr.conf
Environment=PIDFILE=/var/opt/rh/rh-postgresql95/run/repmgr/repmgrd-9.5.pid
