# COPY/etc/systemd/system/rh-repmgr95.service.d/config_file.conf
[Service]
Environment=REPMGRDCONF=/etc/repmgr.conf
