# COPY/usr/lib/systemd/system/evm-failover-monitor.service
[Unit]
Description=PostgreSQL Failover Monitor for ManageIQ
After=evmserverd.service

[Service]
ExecStart=/bin/sh -c /etc/default/evm && bundle exec postgres_ha_admin
Restart=on-failure

[Install]
WantedBy=multi-user.target
