# COPY/usr/lib/systemd/system/evm-failover-monitor.service
[Unit]
Description=PostgreSQL Failover Monitor for ManageIQ
After=manageiq.service

[Service]
ExecStart=/usr/bin/evm-failover-monitor
Restart=on-failure

[Install]
WantedBy=multi-user.target
