# COPY/usr/lib/systemd/system/evm-watchdog.service
[Unit]
Description=evm watchdog process

[Service]
ExecStart=/bin/evm_watchdog
Restart=on-failure

[Install]
WantedBy=multi-user.target
